# tools/dev/docker/supervisord-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startsecs=0
startretries=3
redirect_stderr=true
