# tools/dev/docker/supervisord-onos.conf
[program:onos]
command=/opt/onos/bin/onos-service
autostart=true
autorestart=true
startsecs=0
startretries=3
user=sdn
environment=ONOS_USER="sdn"
redirect_stderr=true
